# .config/systemd/user/eye-strain-notify.service
[Unit]
Description=Prevent Eye Strain.

[Service]
Type=oneshot
ExecStart=%h/.local/bin/eye-strain-notify.sh
