# .config/systemd/user/rclone-bisync.service
[Unit]
Description=Rclone bisync job

[Service]
Type=oneshot
ExecStart=%h/.local/bin/rclone-bisync.sh
